# cron_purge_daily.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2025-09-09T00:00:00.0000000</Date>
    <Author>SAMIA-TAROT</Author>
    <Description>Daily data purge - remove old audit logs and rate limit data (50+ days retention)</Description>
    <URI>\SAMIA-TAROT\Daily Purge</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2025-09-10T00:10:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT10M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>curl</Command>
      <Arguments>-X POST https://<API_BASE_URL>/api/cron/purge_old -H "X-User-ID: <ADMIN_USER_UUID>" -H "X-Job-Token: $env:JOB_TOKEN"</Arguments>
      <WorkingDirectory>C:\</WorkingDirectory>
    </Exec>
  </Actions>
</Task>

<!-- 
SETUP INSTRUCTIONS:
1. Replace <API_BASE_URL> with your API domain (e.g., api.samia-tarot.com)  
2. Replace <ADMIN_USER_UUID> with an actual admin user UUID
3. Set JOB_TOKEN environment variable (generate: openssl rand -hex 32)
4. Import via: schtasks /Create /XML cron_purge_daily.xml /TN "SAMIA-TAROT\Daily Purge"
5. Verify: schtasks /Query /TN "SAMIA-TAROT\Daily Purge"

NOTES:
- Runs at 00:10 UTC daily
- Network connection required
- 10 minute timeout
- Will not start if another instance is running
-->